# supervisor.conf
[program: spring_service4]
command=/home/ubuntu/spring-app/mvnw spring-boot:run
directory=/home/ubuntu/spring-app/service4
environment=DATABASE_URL=jdbc:postgresql://X.X.X.X:5432/database